# config/requirements/conf/supervisor/conf.d/metaverse1.ini
[program:metaverse1]
user=groma
directory=/home/metaverse1/metaverse1/
environment=HOME=/home/groma,
command=/home/metaverse1/metaverse1/manage.py runfcgi
        daemonize=false protocol=fcgi host=127.0.0.1 port=3000
        minspare=1 maxspare=1 maxchildren=1
