# config/requirements/conf/supervisor/conf.d/udp_proxy.ini
[program:udp_proxy]
user=nobody
command=/home/groma/udp_proxy.py
